# systemd/kb-watch.service
[Unit]
Description=kb vault watcher — incremental reindex on save
Documentation=https://github.com/nubank/diraol-personal-kb
After=default.target

[Service]
Type=simple
WorkingDirectory=%h/kb
ExecStart=%h/.local/bin/uv run --project %h/kb kb-watch
Restart=on-failure
RestartSec=5s
StandardOutput=append:%h/.local/state/kb-watch.log
StandardError=append:%h/.local/state/kb-watch.log
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
